# root/etc/inittab
::sysinit:/sbin/openrc sysinit
::sysinit:/sbin/openrc boot
::wait:/sbin/openrc default
::shutdown:/sbin/openrc shutdown

tty1::respawn:/sbin/getty 38400 tty1
ttyS0::respawn:/bin/sh -c 'stty cols 80 rows 24 2>/dev/null; exec btop --force-utf'
